# packaging/linux/root/usr/lib/systemd/user/voxspell.service
[Unit]
Description=VoxSpell voice input daemon
StartLimitIntervalSec=30
StartLimitBurst=5
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/voxspell-daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
